# urxvtd.service
[Unit]
Description=RXVT-Unicode Daemon

[Service]
User=flash
ExecStart=/usr/bin/urxvtd -q -o

[Install]
WantedBy=multi-user.target
